# services/sphinx/run
#!/bin/bash
exec /sbin/setuser abc /usr/bin/searchd --nodetach --config /config/etc/sphinx.conf
